# pxe/config/inittab
# We have to do the other stuff first
::sysinit:/etc/init.d/rcS

#  Let's get some shells
::askfirst:/bin/sh
# We don't have any terminals
#tty1::askfirst:/bin/sh
#tty2::askfirst:/bin/sh
#tty3::askfirst:/bin/sh
#tty4::askfirst:/bin/sh

# Example of how to put a getty on a serial line (for a terminal)
#
::respawn:/sbin/getty -l /bin/sh -L -n ttyS0 9600 vt100


# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
#::shutdown:/sbin/swapoff -a
